# daemon/ya-ytdlp.service
[Unit]
Description=YT-DLP Server for subtitles

[Service]
ExecStart=/usr/bin/node /root/proj/ytdlp/index.js
Restart=always

[Install]
WantedBy=multi-user.target